# scripts/com.timedtrading.calibrate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.timedtrading.calibrate</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd /Users/shashant/timedtrading &amp;&amp; ./scripts/sync-d1.sh &amp;&amp; node scripts/calibrate.js</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>0</integer>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/shashant/timedtrading/data/calibration.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/shashant/timedtrading/data/calibration.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/shashant/timedtrading</string>
</dict>
</plist>
